# server/rfzv45-agent.service
[Unit]
Description=RFZv45 如意 Agent Server
After=network.target surrealrdb.service

[Service]
Type=simple
User=ubuntu
WorkingDirectory=/home/wang/RFZv45/server
Environment=DEEPSEEK_API_KEY=sk-placeholder
Environment=PYTHONUNBUFFERED=1
ExecStart=/usr/bin/python3 agent.py
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
